# startup.conf
description "my app"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec /home/vagrant/app/env/bin/uwsgi /home/vagrant/app/production.ini
